# examples/rc.local
# Add this to /etc/rc.local to autoload when not starting X
(sleep 4; sudo -u pi key-mapper-control --command autoload) &

